# docker/osrm/server/service/routed/run
#!/bin/sh

exec /usr/bin/osrm-routed --port 5000 --shared-memory --threads 8 --max-table-size 100000 2>&1
